# 90-usbmidiconnect.rules
ACTION=="add",\
SUBSYSTEMS=="usb",\
RUN+="/bin/bash /home/pi/ysynth4/midiconnect.sh"